# ansible/roles/prh_tools/templates/prh.timer
[Unit]
Description=Run prh service at midnight

[Timer]
OnCalendar=*-*-* 0:00:00
Persistent=True

[Install]
WantedBy=timers.target
