# ansible/roles/prh_tools/templates/prh_clear.service
[Unit]
Description=Clears PRH data

[Service]
ExecStart=/usr/lib/ckan/default/bin/ckan --config {{ ckan_ini }} prh-tools clear /srv/ytp/files/prh

